# systemd/user/pulseaudio.service.d/override.conf
[Service]
# Have to clear ExecStart first
ExecStart=
ExecStart=/usr/bin/pulseaudio --daemonize=no --log-target=file:/home/max/.cache/pulseaudio/log.txt --log-level=info
