# systemd/user/redshift-gtk.service.d/override.conf
[Unit]
Requires=x11.target
After=x11.target

[Install]
WantedBy=graphical-session.target
